# powerdns/rootfs/etc/sv/pdns/run
#!/bin/sh -e

exec 2>&1
exec chpst -updns pdns_server --daemon=no
