# haraka/rootfs/etc/sv/haraka/run
#!/bin/sh -e

source /etc/envvars

export PATH=$PATH:/usr/local/bin

exec 2>&1
exec /usr/local/bin/haraka -c /app/haraka 2>&1
