# haraka/rootfs/etc/sv/haraka/run
#!/bin/sh -e

source /etc/envvars
export PATH=$PATH:/usr/local/bin
HARAKA=/app/haraka

exec 2>&1
exec /usr/local/bin/haraka -c $HARAKA 2>&1
