# haraka/rootfs/etc/sv/haraka/run
#!/bin/sh -e

source /etc/envvars

cd /app/haraka

# restore backup folder
if [ ! -f /app/haraka/package.json ]; then
  rsync --update -raz --progress --exclude="config" /app/haraka-bak/ /app/haraka

fi

if [ ! -f /app/haraka/config/me ]; then
  rsync --update -raz --progress /app/haraka/config/ /app/haraka/config
fi

echo "$SMTP_HOST" > /app/haraka/config/me
chown -R node:node /app/haraka/haraka

exec 2>&1
exec chpst -u node:node env HOME='/app/haraka' haraka -c /app/haraka/haraka 2>&1
