# mariadb/rootfs/etc/service/mysqld/run
#!/bin/sh -e

exec 2>&1
exec /usr/bin/mysqld --user=mysql --console --pid-file=/var/run/mysqld/mysqld.pid --socket=/var/run/mysqld/mysqld.sock --port=3306
